# modules/ybnet/files/monitor/monitor/src/cron.d/monitor-cron
* * * * * root /usr/local/bin/monitor/start.sh monitor >/dev/null 2>&1
*/5 * * * * root /usr/local/bin/monitor/start.sh game-monitor >/dev/null 2>&1
0 15 * * 1 root kill -9 $(pgrep python3) && nohup python3 /usr/local/bin/monitor/src/mods/plugins/system-exporter.py >/dev/null 2>&1 &
